# host-overrides/jesterpm.net/config/systemd/user/dyndns.timer
[Unit]
Description=DynDNS Service

[Timer]
OnCalendar=*-*-* *:00/5:00
Persistent=true
Unit=dyndns.service

[Install]
WantedBy=timers.target
